# etc/cron.daily/fstrimDaemon
#!/bin/sh

exec /usr/sbin/fstrimDaemon.sh one_shot